# zram-writeback-idle@.timer
[Unit]
Description=Mark pages in /dev/%i as idle
DefaultDependencies=false
ConditionPathExists=/sys/block/%i/idle

[Timer]
OnUnitActiveSec=10min
OnActiveSec=10min
AccuracySec=10s
RandomizedDelaySec=10s
FixedRandomDelay=true

[Install]
WantedBy=timers.target
